# launchd/com.epilog.cloudflare.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.epilog.cloudflare</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>trend-watcher</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/abimangkuagent/EPILOG-TREND-ANALYZER</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Asia/Jakarta</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/abimangkuagent/EPILOG-TREND-ANALYZER/logs/cloudflare-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/abimangkuagent/EPILOG-TREND-ANALYZER/logs/cloudflare-stderr.log</string>
</dict>
</plist>
